# roles/quay_deploy/templates/quay-clair.container.j2
[Unit]
Description=Clair vulnerability scanner for Red Hat Quay
Documentation=https://docs.redhat.com/en/documentation/red_hat_quay/
After=network-online.target quay-network-network.service quay-clair-postgres.service quay.service
Wants=network-online.target
Requires=quay-network-network.service quay-clair-postgres.service

[Container]
Image={{ clair_image }}:{{ clair_image_tag }}
ContainerName=quay-clair
Network=quay-network.network

Environment=CLAIR_MODE=combo
Environment=CLAIR_CONF=/clair/config.yaml

Volume={{ quay_clair_config_dir }}:/clair:Z

PodmanArgs=--cpus={{ quay_clair_cpus }}

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
