# roles/quay_deploy/templates/quay-clair-postgres.container.j2
[Unit]
Description=PostgreSQL database for Clair vulnerability scanner
After=network-online.target quay-network-network.service
Wants=network-online.target
Requires=quay-network-network.service

[Container]
Image={{ postgres_image }}:{{ postgres_image_tag }}
ContainerName=quay-clair-postgres
Network=quay-network.network
Volume=quay-clair-pgdata.volume:/var/lib/pgsql/data:Z
EnvironmentFile={{ quay_quadlet_dir }}/quay.env

Environment=POSTGRESQL_USER=clair
Environment=POSTGRESQL_DATABASE=clair
Environment=POSTGRESQL_MAX_CONNECTIONS={{ quay_clair_postgres_max_connections }}

PodmanArgs=--cpus={{ quay_clair_postgres_cpus }}

HealthCmd=pg_isready -U clair -d clair
HealthInterval=10s
HealthRetries=5
HealthTimeout=5s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
